# consul/etc/service/consul/run
#!/bin/sh
exec consul agent -ui-dir /consul-ui -data-dir /consul-data -config-dir /consul